# config/systemd-ask-password-wrapper.path
[Unit]
Description=Forward Password Requests to Wrapper Directory Watch

DefaultDependencies=no
Before=paths.target cryptsetup.target
Conflicts=emergency.service
Before=emergency.service
Conflicts=shutdown.target
Before=shutdown.target

[Path]
DirectoryNotEmpty=/run/systemd/ask-password
MakeDirectory=yes
